# os/randomness.service
[Unit]
Description=Random of Randomness
After=network.target
Wants=network-online.target

[Service]
User=APP_USER
Group=APP_USER
WorkingDirectory=REPO_DIR

# Environment variables (permissions 600, owned by APP_USER)
EnvironmentFile=REPO_DIR/.env
Type=oneshot
ExecStart=REPO_DIR/scripts/daemon.bash
TimeoutStartSec=300

# No restart needed for a short-lived batch job
Restart=no

# Capture logs in journald
StandardOutput=journal
StandardError=journal

# Optionally increase open file limit
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
